# files/supervisord_cti_api.ini
[program:cti-api]
directory = /usr/local/src/ikono/cti-api
command = /root/.envs/cti-api/bin/python init_cti_api.py
autostart=true
autorestart=true
user=root
stdout_logfile = /usr/local/src/ikono/cti-api/logs/supervisor_out.log           ; Where to write log messages
stderr_logfile = /usr/local/src/ikono/cti-api/logs/supervisor_err.log           ; Where to write log messages
redirect_stderr = false                                                           ; Save stderr in the same log
exitcodes=0,1
